# systemd_daemons/factorio-updater.timer
[Unit]
Description=Timer for Factorio Updater Service

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
